# sdr/sdr-temps.service
[Unit]
Description=SDR Temperature Logger + Web Dashboard
After=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
User=jsn
Environment=HOME=/home/jsn
WorkingDirectory=/home/jsn/cc_projects/sdr

ExecStart=/usr/bin/tmux new-session -d -s sdr \
    "rtl_433 -F json -M time:iso 2>/dev/null | python3 temp_server.py --dir sdr"

ExecStop=/usr/bin/tmux kill-session -t sdr

TimeoutStartSec=30
TimeoutStopSec=15

Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
